# service/stock-monitor.service.template
[Unit]
Description=Stock Monitor (${LABEL})
After=network.target

[Service]
Type=simple
WorkingDirectory=${PROJECT_ROOT}
ExecStart=${VENV_PYTHON} -m stock_monitor --config ${CONFIG_PATH} --host ${HOST} --port ${PORT}
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PATH=${SYSTEM_PATH}

StandardOutput=append:${LOG_DIR}/stdout.log
StandardError=append:${LOG_DIR}/stderr.log

[Install]
WantedBy=default.target
